# systemd/hera-hashpipe-monitor.service
# Configuration file for systemd that keeps the HERA hashpipe monitor running.
#
# Copy this file to /etc/systemd/system/. Then run
# `systemctl enable hera-hashpipe-monitor` and `systemctl start hera-hashpipe-monitor`.
#
# This service is meant to be run on hera-corr-head.

[Unit]
Description=HERA Hashpipe monitor to InfluxDB

[Service]
Type=simple
Restart=always
RestartSec=60
User=hera
Group=hera
ExecStart=/home/hera/hera-venv/bin/hera_hashpipe_monitor.py

[Install]
WantedBy=multi-user.target
